# supervisord.conf
; supervisord.conf — Contexta program definitions
;
; Included by /etc/supervisor/supervisord.conf (the Debian package default).
; Manages two co-equal processes inside the single Docker container:
;
;   api    — uvicorn serving FastAPI on port 8000 (data REST layer)
;   reflex — Reflex serving the compiled UI on port 3000 and handling
;            WebSocket state-sync on port 8001
;
; Environment variables are forwarded from the container environment.

[program:api]
command=uvicorn contexta.api:app --host 0.0.0.0 --port 8000 --workers 1
directory=/app
autostart=true
autorestart=true
; Give the DB a few seconds to be available before accepting requests.
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
    CONTEXTA_DB_PATH="%(ENV_CONTEXTA_DB_PATH)s",
    CONTEXTA_EXPORT_PATH="%(ENV_CONTEXTA_EXPORT_PATH)s",
    CONTEXTA_LOG_LEVEL="%(ENV_CONTEXTA_LOG_LEVEL)s"

[program:reflex]
; --env prod  : serves the pre-built Next.js bundle (no hot-reload)
; --backend-port 8001 : avoids collision with FastAPI on 8000
; --frontend-port 3000 : the port exposed to the end user
command=python -m reflex run --env prod --backend-port 8001 --frontend-port 3000
directory=/app
autostart=true
autorestart=true
; Reflex needs a few extra seconds to warm up the WebSocket event loop.
startsecs=8
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
    CONTEXTA_API_URL="%(ENV_CONTEXTA_API_URL)s"
